# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[include]
files = /tmp/supervisord-*.conf

[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/postgresql.conf
autostart=false
autorestart=true
priority=5
user=postgres
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=false
autorestart=true
startsecs=2
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=secureyeoman
environment=XDG_DATA_HOME="/data/caddy"

[program:llm_gateway]
command=/usr/local/bin/llm_gateway daemon
autostart=true
autorestart=true
startsecs=2
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=secureyeoman

[program:agent_runtime]
command=/usr/local/bin/agent_runtime daemon
autostart=true
autorestart=true
startsecs=2
startretries=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=secureyeoman

[program:secureyeoman]
command=/bin/sh -c "exec %(ENV_SY_CMD)s"
autostart=true
autorestart=true
startsecs=5
startretries=10
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=secureyeoman
